# tools/98-huawei_modem.rules
ATTRS{idVendor}=="12d1", ATTRS{idProduct}=="1506", RUN+="/tools/modem/modem_start.sh"
